# supervisor/gitlab.conf
[program:gitlab]
directory=/home/git/gitlab
command=/home/git/gitlab/bin/web start_foreground
user=git
