# supervisor/sidekiq.conf
[program:sidekiq]
directory=/home/git/gitlab
command=/home/git/gitlab/bin/background_jobs start_no_deamonize
user=git
